# deploy/cron/crontab
* * * * * php /var/www/html/cep_backend/artisan schedule:run >> /dev/null 2>&1